# launchd/com.goodidea.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.goodidea.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>SCRIPT_DIR/sync/sync_goodidea.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>SCRIPT_DIR</string>

  <!-- Every Monday at 09:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ANTHROPIC_API_KEY</key>
    <string>REPLACE_WITH_YOUR_KEY</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>SCRIPT_DIR/logs/sync_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>SCRIPT_DIR/logs/sync_stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
